# bull/fussbull.service
[Unit]
Description=Fussbull
[Service]
ExecStart=/usr/bin/fussbull
Type=oneshot
[Install]
WantedBy=multi-user.target
